# roles/hadoop/templates/etc/profile.d/hadoop.sh

export HADOOP_CONF_DIR=/usr/src/hadoop-3.0.0/etc/hadoop/
export YARN_CONF_DIR=/usr/src/hadoop-3.0.0/etc/hadoop/
export HADOOP_USER_NAME=hadoop
export HADOOP_HOME=/usr/src/hadoop-3.0.0
export PATH=$PATH:$HADOOP_HOME/bin
