# roles/tez/templates/etc/profile.d/tez.sh
# {{ ansible_managed }}

export TEZ_JARS=/usr/src/apache-tez-0.9.1
export TEZ_CONF_DIR=/usr/src/apache-tez-0.9.1/conf/
# export HADOOP_CLASSPATH=$HADOOP_CLASSPATH:${TEZ_CONF_DIR}:${TEZ_JARS}/*:${TEZ_JARS}/lib/*
